# assets/rtorrent-supervisord.ini
[program:init]
autostart = true
autorestart = false
startsecs = 0
user = root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
command = /bin/sh /usr/local/bin/init.sh

[program:rtorrent]
autostart = false
autorestart = true
user = root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
command = /bin/sh /usr/local/bin/start-rtorrent.sh

[program:lighttpd]
autostart = false
autorestart = true
user = root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
command = /bin/sh /usr/local/bin/start-lighttpd.sh